# docker/cloudwatch/amazon-cloud-watch.conf
[program:amazon-cloud-watch]
command=/opt/aws/amazon-cloudwatch-agent/bin/start-amazon-cloudwatch-agent
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stopwaitsecs=3600
